# Launchpad/uhub-unstable/debian/uhub-unstable.postinst
#!/bin/sh

set -e

case "$1" in
    configure)

    chmod 0750 /var/log/uhub

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst: error: unknown argument: $1" >&2
        exit 1
        ;;
esac

#
if [ -x "/etc/init.d/uhub" ]; then
	update-rc.d uhub defaults >/dev/null
	invoke-rc.d uhub start || exit $?
fi
#

#DEBHELPER#

exit 0
